# subsonic-installer-debian/src/DEBIAN/preinst
#! /bin/sh

set -e

# Stop Subsonic service.
if [ -e /etc/init.d/subsonic ]; then
  invoke-rc.d subsonic stop
fi

# Backup database.
if [ -e /var/subrobotnik/db ]; then
  rm -rf /var/subrobotnik/db.backup
  cp -R /var/subrobotnik/db /var/subrobotnik/db.backup
fi

